# supervisor_tomcat.conf
[program:tomcat]
directory=/
command=/usr/local/tomcat/bin/catalina.sh run
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log